# conf/supervisord.conf
#!/bin/bash

[supervisord]
nodaemon=true

[program:apache2]
command=/bin/bash -c "/usr/sbin/httpd -D  FOREGROUND"

[program:catalina]
command=/bin/bash -c "/usr/libexec/tomcat/server start"




